# debian/geronimo.cron
# Die untenstehenden Imports benoetigen jeweils einen Datenbank-Lock - also nie parallel ausfuehren
# Wiki: stuendlich
13 * * * *      geronimo-manage import_wiki >/dev/null
# ondata: die Verfallsdauer der Informationen ist aufgrund eines olsr-ondataservice-Bugs wohl recht kurz (Minuten)
*/15 * * * *    geronimo-manage import_olsr http://localhost:2006/; geronimo-manage import_ondataservice
