# debian/on-freifunk-api.cron.d
*/10 *  * * *   on-freifunk-api    /usr/bin/on-freifunk-api-update-nodelist >/var/lib/on-freifunk-api/public/api.freifunk.net-nodelist.json
